# srvfalar/etc/systemd/system/srvfalar.service
[Unit]
Description=srvfalar

[Service]
ExecStart=/usr/bin/srvfalar

[Install]
WantedBy=multi-user.target
